# supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

#redis
[program:redis]
command=redis-server --save "" --appendonly no
autostart=true
autorestart=true
priority=1

[program:gpt_queue_1]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_2]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_3]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_4]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_5]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_6]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_7]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_8]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_9]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gpt_queue_10]
command=python manage.py rqworker gpt
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:gemini_queue_1]
command=python manage.py rqworker gemini
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gemini_queue_2]
command=python manage.py rqworker gemini
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
# celery workers
[program:celery_queue1]
command=celery -A upscale.celery.app worker -Q delay_queue --autoscale=15,5 --loglevel=info
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:celery_queue2] 
command=celery -A upscale.celery.app worker -Q risk_queue --autoscale=15,5 --loglevel=info
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:celery_queue3]
command=celery -A upscale.celery.app worker -Q report_queue --autoscale=15,5 --loglevel=info
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0


[program:celery_queue4]
command=celery -A upscale.celery.app worker -Q sch_queue --autoscale=15,5 --loglevel=info
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:celery_queue5] 
command=celery -A upscale.celery.app worker -Q drone_queue --autoscale=15,5 --loglevel=info
autostart=true
autorestart=true
priority=2
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0



#app
[program:app]
command=python run_app.py
directory=/app
autostart=true
autorestart=true
priority=3
startsecs=5
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

